# custom_service.service.in
[Unit]
Description=@PROJECT_DESCRIPTION@

[Service]
Type=simple
WorkingDirectory=@CMAKE_INSTALL_PREFIX@
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/@PROJECT_NAME@ user=%u hostname=%H (boot=%b)
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
